# system/kiosk.service
[Unit]
Description=Smart Kiosk UI
After=network-online.target

[Service]
Type=simple
Environment=SMART_KIOSK_CONFIG=/etc/smart-kiosk/config.yaml
ExecStart=/usr/bin/python3 -m smart_kiosk.cli run
Restart=always
RestartSec=3
WorkingDirectory=/opt/smart-kiosk

[Install]
WantedBy=multi-user.target

